# crontab
MAILTO=postmaster@digitalkingdom.org
LANG=en_US.UTF-8
3 3 * * * /home/sampre_mw/mediawiki/cron_backup.sh
